# .config/systemd/system/asus-fans-on.service
[Unit]
Description=ASUS Fans On
[Service]
ExecStart=/usr/lib/systemd/system/asus-fans-on.sh
[Install]
WantedBy=multi-user.target
