# .config/systemd/system/eno1-1000.service
[Install]
WantedBy=multi-user.target

[Unit]
Description=ethtool eno1 gigabit

[Service]
ExecStart=ethtool -s eno1 speed 1000
